# recipes-kernel/linux-firmware/files/monitor-m4-elf-file.path
[Unit]
Description=Monitor /tmp/arduino/m4-user-sketch.elf file
DefaultDependencies=no
Requires=m4-led-forwarder.service
After=m4-led-forwarder.service

[Path] 
PathModified=/tmp/arduino/m4-user-sketch.elf

[Install]
WantedBy=basic.target
