# recipes-support/arduino-ootb/arduino-ootb/connection-status-led.timer
[Unit]
Description=Timer connection status led

[Timer]
OnCalendar=*:*:0/5
Persistent=true

[Install]
WantedBy=timers.target

