# recipes-support/arduino-ootb/arduino-ootb/connection-status-led.service
[Unit]
Description=Connection status led
After=network-online.target

[Service]
LogLevelMax=3
Type=oneshot
ExecStart=/bin/sh /usr/bin/connection_status_led
