# recipes-devtools/android-tools/android-tools/android-tools-adbd.service
[Unit]
Description=Android Debug Bridge
After=network-pre.target dbus.service
Before=android-system.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/arduino-android-gadget-setup

[Install]
WantedBy=basic.target
